# hosts/cltc-aus-lws03/.config/systemd/user/dotfiles-autotheme.service.d/hyprland-only.conf
# cltc-aus-lws03 is a dual GNOME+Hyprland (GDM) Fedora box; graphical-session.target is reached
# under both. The wallpaper watcher polls hyprpaper, so it must only run under
# Hyprland. Depends on the session env being imported into the systemd --user
# manager, which hosts/cltc-aus-lws03/.config/hypr/local.conf does
# (dbus-update-activation-environment) on Hyprland start.
[Unit]
ConditionEnvironment=XDG_CURRENT_DESKTOP=Hyprland
