# hosts/fedora/.config/systemd/user/waybar.service.d/hyprland-only.conf
# Fedora's packaged waybar.service is WantedBy=graphical-session.target, which is
# also reached under GNOME. Gate it to the Hyprland session only (see the kanshi
# drop-in and hosts/fedora/.config/hypr/local.conf for the env import it needs).
[Unit]
ConditionEnvironment=XDG_CURRENT_DESKTOP=Hyprland
