# profiles/hyprland/.config/systemd/user/waybar.service.d/hyprland-only.conf
# Drop-in for the Fedora-packaged waybar.service
# (/usr/lib/systemd/user/waybar.service).
#
# graphical-session.target is also active under GNOME on this dual-session
# machine, and the stock unit only conditions on WAYLAND_DISPLAY (which is set
# under GNOME Wayland too). Without this, enabling waybar.service would put
# Waybar on top of the GNOME session. Gate it to the Hyprland session instead.
#
# NOTE: ConditionEnvironment is ANDed with the stock unit's conditions, so this
# only narrows when Waybar starts; it never widens it. Relies on hyprland.conf
# importing the session environment into the systemd --user manager on start.
[Unit]
ConditionEnvironment=XDG_CURRENT_DESKTOP=Hyprland
