# profiles/hyprland/.config/systemd/user/elephant.service
[Unit]
Description=Elephant data provider backend (walker)
Documentation=https://github.com/abenz1267/elephant
PartOf=graphical-session.target
After=graphical-session.target
# walker's backend; only needed under Hyprland on this dual-session box.
# (Same guard pattern as kanshi.service / the waybar drop-in.)
ConditionEnvironment=XDG_CURRENT_DESKTOP=Hyprland

[Service]
Type=simple
# Absolute path via %h: the binary lives in ~/.local/bin, which the systemd
# user manager does NOT search for bare ExecStart names (unlike Hyprland's PATH).
ExecStart=%h/.local/bin/elephant
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
