# profiles/hyprland/.config/systemd/user/walker.service
[Unit]
Description=Walker launcher (GApplication service, walker's fast-open backend)
Documentation=https://github.com/abenz1267/walker
PartOf=graphical-session.target
After=graphical-session.target elephant.service
# walker needs elephant's socket; order after it and pull it in.
Wants=elephant.service
# Hyprland-only on this dual-session box (same guard as elephant/kanshi/waybar).
ConditionEnvironment=XDG_CURRENT_DESKTOP=Hyprland

[Service]
Type=simple
ExecStart=%h/.local/bin/walker --gapplication-service
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
